# etc/profile.d/10-android-env
# Set the SDK directory based off of android-studio's location
[ -d "${ANDROID_HOME:=${HOME}/Library/Android/sdk}" ] && export ANDROID_HOME && {
	source update-env post PATH "${ANDROID_HOME}/tools"
	source update-env post PATH "${ANDROID_HOME}/tools/bin"
	source update-env post PATH "${ANDROID_HOME}/platform-tools"

	[ -d "${ANDROID_NDK_HOME:=${ANDROID_HOME}/ndk-bundle}" ] && export ANDROID_NDK_HOME && {
		source update-env post PATH "${ANDROID_NDK_HOME}"
	}|| {
		unset ANDROID_NDK_HOME
	}	
} || {
    unset ANDROID_HOME
}
