# etc/profile.d/10-bash-completion-env
if [ "${BASH-no}" != "no" -a -f $(brew --prefix)/etc/bash_completion ]; then
  . $(brew --prefix)/etc/bash_completion
fi
