# etc/profile.d/10-csshx-env
[ -z "$BASH_VERSION" -o -z "$PS1" -o "${BASH-no}" != "/bin/bash" ] && return
which csshX &>/dev/null && {
  eval "$(csshX --bash_env)"
}
