# etc/profile.d/10-java-env
if [ -x /usr/libexec/java_home ]; then
    export JAVA_HOME="$(/usr/libexec/java_home -v1.8 2>/dev/null)"
    
    [ -z "$BASH_VERSION" -o -z "$PS1" -o "${BASH-no}" != "/bin/bash" ] && return
	declare -a ON_JAVA_CHANGE
	use-java() { 
		export JAVA_HOME="$(/usr/libexec/java_home -v1.${1} 2>/dev/null)"
        for i in "${ON_JAVA_CHANGE[@]}"; do ${i}; done
	}
fi
