# etc/profile.d/10-java-env
if [ -x /usr/libexec/java_home ]; then
    export JAVA_HOME="$(/usr/libexec/java_home -v1.8 2>/dev/null)"
fi
