# etc/profile.d/10-maven-env
which mvn >/dev/null && {
    export MAVEN_OPTS="${MAVEN_OPTS} -Xmx512m -Xms512m -XX:MaxPermSize=256m -Djava.awt.headless=true"
}
