# etc/profile.d/10-maven-env
which mvn >/dev/null && {
	[ -z "${M2_HOME}" ] && { M2_HOME="$(brew --prefix maven 2>/dev/null)/libexec" && export M2_HOME; }
    source update-env post MAVEN_OPTS "-Xmx512m -Xms512m -XX:MaxPermSize=256m -Djava.awt.headless=true" " "
}
