# etc/profile.d/10-perl-local-env
PERL_BASE="${HOME}/.perl5"
source update-env pre PATH "${PERL_BASE}/bin" ":"
source update-env pre PERL5LIB "${PERL_BASE}/lib/perl5" ":"
source update-env pre PERL_LOCAL_LIB_ROOT "${PERL_BASE}" ":"
export PERL_MB_OPT="--install_base \"${PERL_BASE}\""
export PERL_MM_OPT="INSTALL_BASE=${PERL_BASE}"
