# etc/profile.d/20-google-cloud-sdk-env
GCLOUD_DIR="$(brew --repository)/Caskroom/google-cloud-sdk/latest/google-cloud-sdk"
[ -d "${GCLOUD_DIR}" ] && {
    export GCLOUD_DIR
    source "${GCLOUD_DIR}/path.bash.inc"
    source "${GCLOUD_DIR}/completion.bash.inc"
}
