# etc/profile.d/90-ssh-rmate-env
if [ -n "${SSH_CONNECTION}" ]; then
    which mate &>/dev/null && which rmate &>/dev/null && {
        alias mate=rmate
    }
fi
